# launchd/com.longfeiguo.dingtalk.study.evening.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.longfeiguo.dingtalk.study.evening</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/longfeiguo/.codex/dingtalk-study/send_dingtalk_text.sh</string>
    <string>/Users/longfeiguo/.codex/dingtalk-study/dingtalk.env</string>
    <string>evening</string>
  </array>
  <key>RunAtLoad</key>
  <false/>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>22</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/longfeiguo/.codex/dingtalk-study/logs/dingtalk-study-evening.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/longfeiguo/.codex/dingtalk-study/logs/dingtalk-study-evening.log</string>
  <key>WorkingDirectory</key>
  <string>/Users/longfeiguo/.codex/dingtalk-study</string>
</dict>
</plist>
